# .files/profile.d/macos.sh
# shellcheck shell=bash

[ -d "/Applications/VMware Fusion.app/Contents/Library" ] && PATH="$PATH:/Applications/VMware Fusion.app/Contents/Library"
