# .files/profile.d/pgsql
# shellcheck shell=bash

[ -d /usr/pgsql-9.4/bin ] && PATH=$PATH:/usr/pgsql-9.4/bin
